# src/big-data-stack-chameleon/roles/supervisord/templates/program.ini.j2
[program:{{ item.key }}]
{% for key, value in item.value.iteritems() %}
{{ key }}={{ value }}
{% endfor %}
